# services/mpdscribble.service
[Unit]
Description=MPD Scrobbler
After=mpd.service

[Service]
ExecStart=/usr/bin/mpdscribble --no-daemon

[Install]
WantedBy=default.target
