# services/mpdscribble.service
[Unit]
Description=MPD Scrobbler

[Service]
ExecStart=/usr/bin/mpdscribble --no-daemon

[Install]
WantedBy=multi-user.target
