# services/xcmenu.service
[Unit]
Description=xcmenu clipboard manager

[Service]
ExecStart=/usr/bin/xcmenu
Environment=DISPLAY=:0

[Install]
WantedBy=default.target
